# govtchat.service
# /etc/systemd/system/govtchat.service
#
# Systemd unit file for the GovOps Chat Agent API Service
# Run as user 'vpa' with Conda environment 'cogops'.

[Unit]
Description=GovOps Chat Agent API Service
After=network.target redis.service jiggasha-search.service
Requires=jiggasha-search.service

[Service]
Type=simple
User=vpa
Group=vpa

WorkingDirectory=/home/vpa/CogOpsCB

EnvironmentFile=/home/vpa/CogOpsCB/.env
Environment=HOME=/home/vpa
Environment=VIRTUAL_ENV=/home/vpa/miniconda3
Environment=PATH=/home/vpa/miniconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/vpa/CogOpsCB

ExecStart=/home/vpa/miniconda3/bin/uvicorn api:app --host 0.0.0.0 --port 9000 --workers 2 --timeout-keep-alive 300

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=govtchat-api

[Install]
WantedBy=multi-user.target
